# TLSGateNGv4-main 4/config-files/System-Services/Services/tlsgateNGv6.service
# =============================================================================
# TLSGate NG v6 - IPv6 Worker Service
# =============================================================================
# Startreihenfolge:
#   1. tlsgateNG-poolgen-primary.service
#   2. tlsgateNG-poolgen-backup.service
#   3. tlsgateNGv4.service
#   4. tlsgateNGv6.service                ← DIESER SERVICE
#
# WICHTIG: Startet ERST wenn Poolgen läuft!
#
# Installation:
#   cp tlsgateNGv6.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable tlsgateNGv6.service
#
# =============================================================================

[Unit]
Description=TLSGate NG v4.36 - IPv6 TLS Termination Worker
Documentation=https://github.com/TorstenJahnke/TLSGateNGv4

# Netzwerk muss verfügbar sein
After=network.target network-online.target
Wants=network-online.target

# WICHTIG: Wartet auf Poolgen!
After=tlsgateNG-poolgen-primary.service tlsgateNG-poolgen-backup.service
Requires=tlsgateNG-poolgen-primary.service

# Teil des TLSGate Targets
PartOf=tlsgateNG.target

[Service]
Type=simple

# Worker läuft als unprivilegierter User
User=tlsgateNG
Group=tlsgateNG

# Arbeitsverzeichnis
WorkingDirectory=/etc/tlsgateNG

# IPv6 Worker mit SHM-Keypool
ExecStart=/usr/local/bin/tlsgateNGv6 \
    --shm \
    -c /etc/tlsgateNG/tlsgateNG.conf

# Restart bei Crash
Restart=always
RestartSec=3s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

# Performance für High-Traffic
LimitNOFILE=1048576
LimitNPROC=4096
Nice=-5

# Security Hardening (Worker hat nur lesenden SHM-Zugriff)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/tlsgateNG /var/log/tlsgateNG /opt/TLSGateNX

# Capabilities (nur Netzwerk)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tlsgateNGv6

[Install]
WantedBy=multi-user.target tlsgateNG.target
