# TLSGateNGv4-main 4/config-files/System-Services/Services/tlsgateNG-poolgen-backup.service
[Unit]
Description=TLSGate NG v4.36 GEN4 - Keypool Generator BACKUP (HA Standby)
Documentation=https://github.com/TorstenJahnke/TLSGateNGv4
After=network.target network-online.target tlsgateng-poolgen-primary.service
Wants=network-online.target
PartOf=tlsgateNG.target

# HA Configuration:
# This is the BACKUP poolgen instance.
# - Starts in STANDBY mode (waiting for lock)
# - If primary crashes or stops: automatically takes over
# - When primary restarts: continues running (first to lock wins)

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/etc/tlsgateNG/poolgen

# Binary path - Keypool generator (BACKUP)
# Config: /etc/tlsgateNG/poolgen/poolgen.conf (copy from examples/)
ExecStart=/usr/local/etc/tlsgateNG/tlsgateNG-poolgen/tlsgateNG-poolgen \
    --poolkeygen \
    --shm \
    --ha-role=backup \
    -c /etc/tlsgateNG/poolgen/poolgen.conf

# Process management
Restart=always
RestartSec=10s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

# Performance tuning
LimitNOFILE=65536
LimitNPROC=512
Nice=-10

# Security hardening (relaxed for shared memory access)
PrivateTmp=no
PrivateIPC=no
NoNewPrivileges=no
ProtectSystem=off
ProtectHome=no

# Capabilities (needed for shared memory and lock file)
AmbientCapabilities=CAP_IPC_LOCK CAP_SYS_RESOURCE
CapabilityBoundingSet=CAP_IPC_LOCK CAP_SYS_RESOURCE CAP_SETUID CAP_SETGID

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tlsgateNG-poolgen-backup

[Install]
WantedBy=multi-user.target
