# TLSGateNGv4-main 4/config-files/System-Services/Services/tlsgateNG-termination01.service
# =============================================================================
# TLSGate NG - Aviontex Termination Worker 01
# =============================================================================
# Startreihenfolge:
#   1. tlsgateNG-poolgen.service
#   2. tlsgateNG-termination01.service  ← DIESER SERVICE
#
# WICHTIG: Startet ERST wenn Poolgen läuft!
#
# Installation:
#   cp tlsgateNG-termination01.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable tlsgateNG-termination01.service
#
# =============================================================================

[Unit]
Description=TLSGate NG v4.36 - Aviontex Termination Worker 01
Documentation=https://github.com/TorstenJahnke/TLSGateNGv4

# Netzwerk muss verfügbar sein
After=network.target network-online.target
Wants=network-online.target

# WICHTIG: Wartet auf Poolgen!
After=tlsgateNG-poolgen.service
Requires=tlsgateNG-poolgen.service

# Teil des TLSGate Targets
PartOf=tlsgateNG.target

[Service]
Type=simple

# Worker läuft als unprivilegierter User
User=tlsgateNG
Group=tlsgateNG

# Arbeitsverzeichnis
WorkingDirectory=/etc/tlsgateNG/aviontex/termination

# 60 Sekunden warten bis Poolgen SHM vollständig initialisiert hat
ExecStartPre=/bin/sleep 60

# Termination Worker mit SHM-Keypool
ExecStart=/usr/local/bin/tlsgateNGv4 \
    --shm \
    -c /etc/tlsgateNG/aviontex/termination/termination.v4.31.conf

# Restart bei Crash
Restart=always
RestartSec=3s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

# Performance für High-Traffic
LimitNOFILE=1048576
LimitNPROC=4096
Nice=-5

# Security Hardening (Worker hat nur lesenden SHM-Zugriff)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/tlsgateNG /var/log/tlsgateNG /opt/TLSGateNX

# Capabilities (nur Netzwerk)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tlsgateNG-termination01

[Install]
WantedBy=multi-user.target tlsgateNG.target
